# .config/systemd/user/zinit-update.service
[Unit]
Description=zinit update

[Service]
Type=oneshot
ExecStart=zsh --singlecommand -c 'zinit self-update' 
